# supervisor/ini/elk.ini
[program:elk-es]
command=/bin/bash -c 'cd /data/elk/elasticsearch/bin/;./elasticsearch -p pid'
user=elk
directory=/data/elk/elasticsearch/bin/
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/supervisor/logs/elk-es.log
